# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/apache2-error.log
stdout_logfile=/var/log/supervisor/apache2-access.log
user=root
priority=10

[program:sendmail]
command=/usr/sbin/sendmail -bd
autostart=%(ENV_USE_LOCAL_MTA)s
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/sendmail-error.log
stdout_logfile=/var/log/supervisor/sendmail-access.log
user=root
priority=20
environment=USE_LOCAL_MTA=false

[include]
files = /etc/supervisor/conf.d/*.conf
